# SystemD/v3/web@.socket
[Unit]
Description=Socket for worker %i

[Socket]
ListenStream = 8888
ReusePort = true
Service = web@%i.service

[Install]
WantedBy = sockets.target
