# SystemD/v2/myservice.socket
[Unit]
Description=My Service Socket
PartOf=myservice@.service

[Socket]
ListenStream=/tmp/mysocket_v2.sock

[Install]
WantedBy=sockets.target