# SystemD/v1/myservice.socket
[Unit]
Description=My Service Socket

[Socket]
ListenStream=/tmp/mysocket_v1.sock
Accept=yes

[Install]
WantedBy=sockets.target